# deploy/gorouter.service
[Unit]
Description=OpenCode Go Aggregator
Documentation=https://github.com/dbowm91/gorouter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gorouter
Group=gorouter
WorkingDirectory=/var/lib/gorouter
ExecStart=/opt/gorouter/.venv/bin/go-aggregator serve --config /etc/gorouter/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# Graceful shutdown
TimeoutStopSec=30
KillSignal=SIGTERM

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/gorouter
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Network
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# System call filtering
SystemCallFilter=@system-service
SystemCallArchitectures=native

# Environment
EnvironmentFile=/etc/gorouter/env

[Install]
WantedBy=multi-user.target
